# meta/recipes-connectivity/neard/neard/neard.service.in
[Unit]
Description=NFC service
After=syslog.target

[Service]
Type=dbus
BusName=org.neard
Restart=on-failure
ExecStart=@installpath@/neard -n
StandardOutput=null

[Install]
WantedBy=multi-user.target
