# meta/recipes-connectivity/neard/neard/neard.service.in
[Unit]
Description=neard service
Documentation=man:neard(8)
After=syslog.target

[Service]
Type=dbus
BusName=org.neard
ExecStart=@installpath@/neard -n

[Install]
Alias=dbus-org.neard.service
WantedBy=multi-user.target
