# meta/recipes-connectivity/gsm/files/gsmd
#! /bin/sh
#
# gsmd  This shell script starts and stops gsmd.
#
# chkconfig: 345 90 40
# description: Gsmd manages access to a serial- or USB-connected GSM
# processname: gsmd

PATH=/bin:/usr/bin:/sbin:/usr/sbin

[ -f /etc/default/rcS ] && . /etc/default/rcS
[ -f /etc/default/gsmd ] && . /etc/default/gsmd

case "$1" in
    start)
        [ -n "$GSM_POW" ] && ( echo "0" >$GSM_POW; sleep 1 )
        [ -n "$GSM_POW" ] && ( echo "1" >$GSM_POW; sleep 1 )
        [ -n "$GSM_RES" ] && ( echo "1" >$GSM_RES; sleep 1 )
        [ -n "$GSM_RES" ] && ( echo "0" >$GSM_RES; sleep 2 )

        echo -n "Starting GSM daemon: "
        start-stop-daemon -S -x /usr/sbin/gsmd -- gsmd -p $GSM_DEV $GSMD_OPTS -d -l syslog

        if [ $? = 0 ]; then
            echo "gsmd."
        else
            echo "(failed.)"
        fi
        ;;
    stop)
        [ -n "$GSM_POW" ] && echo "0" >$GSM_POW

        echo -n "Stopping GSM daemon: "
        start-stop-daemon -K -x /usr/sbin/gsmd
        echo "gsmd."
        ;;
    restart|force-reload)
        $0 stop
        $0 start
        ;;
    *)
        echo "Usage: /etc/init.d/gsmd {start|stop|restart|force-reload}"
        exit 1
        ;;
esac

exit 0
