# meta/packages/networkmanager/files/NetworkManager
#!/bin/sh
#
### BEGIN INIT INFO
# Provides:          NetworkManager
# Required-Start:    $remote_fs dbus hal
# Required-Stop:     $remote_fs dbus hal
# Should-Start:	     $syslog
# Should-Stop:       $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: network connection manager
# Description:       Daemon for automatically switching network 
#		     connections to the best available connection.
### END INIT INFO

. /etc/profile

case $1 in
  'start')
    echo -n "Starting NetworkManager daemon: NetworkManager"
    /usr/sbin/NetworkManager 
    /usr/sbin/NetworkManagerDispatcher
    /usr/sbin/wpa_supplicant -u &
    echo "."
  ;;

  'stop')
    echo -n "Stopping NetworkManager daemon: NetworkManager"
    kill `ps |grep  /usr/sbin/NetworkManagerDispatcher | grep -v grep | cut "-d " -f2`
    kill `ps |grep  /usr/sbin/NetworkManager | grep -v grep | cut "-d " -f2`
    kill `ps |grep  /usr/sbin/wpa_supplicant | grep -v grep | cut "-d " -f2`
    echo "."
  ;;

  'restart')
    $0 stop
    $0 start
  ;;

  *)
    echo "Usage: $0 { start | stop | restart }"
  ;;
esac
