# meta/recipes-core/udev/udev-extraconf/systemd-udevd.service
.include @systemd_unitdir@/system/systemd-udevd.service
[Service]
MountFlags=shared
